# daemon/relay_daemon.conf
description "VTGS Remote Relay Daemon Service"
author  "Zach Leffke <zleffke@vt.edu>"

start on runlevel [234]
stop on runlevel [0156]

chdir /daemon/relay
exec /daemon/relay/relay_daemon.py
respawn
